# parrhasius.service.erb
[Unit]
Description=Run Parrhasius server

[Service]
ExecStart=<%= `which nix-shell`.strip %> --run 'rake run'
EnvironmentFile=<%= ENV['HOME'] %>/parrhasius.env
WorkingDirectory=<%= `pwd`.strip %>

[Install]
WantedBy=multi-user.target
